# server/sodalite.service
[Unit]
Description=Sodalite Media Downloader Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=sodalite
Group=sodalite
WorkingDirectory=/opt/sodalite/server
Environment=PATH=/opt/sodalite/venv/bin
ExecStart=/opt/sodalite/venv/bin/python /opt/sodalite/server/run.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
NoNewPrivileges=yes

# Hardening & Resource Limits
# For disk usage, PrivateTmp creates a private /tmp directory for the service.
# To limit its size to 15GB, you can use: TemporaryFileSystem=/tmp:size=15G
# Or ensure your system's /tmp partition is sized appropriately.
MemoryMax=2G
TasksMax=512
LimitNOFILE=65536
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictRealtime=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sodalite

[Install]
WantedBy=multi-user.target
